# rootfs/etc/services.d/mcp-ha-extended/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant MCP Extended Server - Finish Script
# ==============================================================================

bashio::log.info "Stopping MCP HA Extended Server..."
